# aquamacs/build/org.aquamacs.nightly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>org.aquamacs.nightly</string>
    <key>ProgramArguments</key>
    <array>
        <string>nightly-build.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Minute</key>
        <integer>5</integer>
        <key>Hour</key>
        <integer>1</integer>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Users/dr/Nightly/master</string>
</dict>
</plist>
